# s3-backup.service
[Unit]
Description=back up files to s3

[Service]
Type=oneshot
ExecStart=/usr/local/lib/s3_backup/s3_backup
